# services/tradingbot-alert@.service
[Unit]
Description=Send email alert for failed service %i

[Service]
Type=oneshot
ExecStart=/usr/local/bin/tradingbot-alert.sh %i
